# muneem-frontend.service
[Unit]
Description=MUNEEM Frontend (Static)
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/digi-bahi-ink/dist
ExecStart=/usr/bin/python3 -m http.server 5173
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
